# vbsd/init/vbsd.service.in
[Unit]
Description=Daemon to relay subtitles for VBS
After=syslog.target

[Service]
ExecStart=@SBINDIR@/vbsd
Type=forking
Restart=restart-always

[Install]
WantedBy=multi-user.target

